# systemd/teleglas-alert.service
[Unit]
Description=TELEGLAS Alert Bot (WebSocket)
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/TELEGLAS
Environment=PATH=/opt/TELEGLAS/venv/bin
ExecStart=/opt/TELEGLAS/venv/bin/python -m ws_alert.alert_runner
Restart=always
RestartSec=10
StandardOutput=append:/opt/TELEGLAS/logs/alert-bot.log
StandardError=append:/opt/TELEGLAS/logs/alert-bot.log

# Resource Limits
MemoryMax=256M
CPUQuota=30%

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
